# system/tonado.service
[Unit]
Description=Tonado Music Box
After=network.target mpd.service
Wants=mpd.service

[Service]
Type=simple
User=tonado
Group=tonado
WorkingDirectory=/opt/tonado
ExecStart=/opt/tonado/.venv/bin/uvicorn core.main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
Environment=TONADO_DB_PATH=/opt/tonado/config/tonado.db
Environment=TONADO_MEDIA_DIR=/home/tonado/media
Environment=TONADO_HARDWARE_MODE=auto

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/tonado/config /home/tonado/media /tmp
ProtectHome=true
PrivateTmp=true

# Allow hardware access
SupplementaryGroups=spi i2c gpio audio

[Install]
WantedBy=multi-user.target
